# cron.d/generateGravityWildcards
45 5 * * * root PATH="$PATH:/usr/local/bin/" generateGravityWildcards.sh
