# cron.d/generateGravityWildcards
55 4 * * * root PATH="$PATH:/usr/local/bin/" generateGravityWildcards.sh
